# rootfs/etc/rc.local
#!/bin/bash

sudo nohup /opt/csocso/run.sh >/dev/null &

exit 0
